# etc/init.d/steam_renice_daemon
#!/sbin/openrc-run

pidfile="/run/steam_renice_daemon.${instance}.pid"
command="/usr/local/bin/steam_renice_daemon.sh"
command_background="true"

